# supervisor.d/zoomreports.ini
[program:zoomreports]
command=ash -c "rm -f %(ENV_PID_FILE)s && bundle exec puma"
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
